# cowrie/supervisord.conf
[supervisord]
nodaemon=true

[program:kippo]
user=tpot
directory=/opt/cowrie/base/
command=/opt/cowrie/base/start.sh
autorestart=true
